# user/.config/systemd/user/gpg-agent.service.d/override.conf
[Service]
ExecStart=
ExecStart=gpg-agent --supervised --homedir %E/gnupg
