# user/.config/systemd/user/mpd-notification.service.d/override.conf
[Unit]
After=
After=mpd.service network.target network-online.target
After=dunst.service


[Service]
ExecStart=
ExecStart=/usr/bin/mpd-notification -m "%h/Music" -t 8