# user/.config/systemd/user/mpc-update.service
[Unit]
Description=Update mpd music library
PartOf=graphical-session.target
After=mpd.service

[Service]
ExecStart=mpc update
