# user/.config/systemd/user/systemd-dump-env-at-login.service
[Unit]
Description=Record ENVIRONMENT at startup
PartOf=default.target

[Service]
Type=oneshot
ExecStart=env

[Install]
WantedBy=default.target
