# user/.config/systemd/user/xrdb-load.service
[Unit]
Description=Load Xresources
PartOf=graphical-session.target

[Service]
ExecStart=xrdb -load "$XDG_CONFIG_HOME/X11/.Xresources"
